# services/tailscale-port-dashboard/ops/tailscale-port-dashboard.user.service
[Unit]
Description=Tailscale port dashboard user service example
After=default.target

[Service]
Type=simple
WorkingDirectory=/opt/microservices/services/tailscale-port-dashboard
Environment=PORT=80
ExecStart=/usr/local/bin/bun run start
Restart=on-failure
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=default.target
